# webapp/xitoy_webapp.service
[Unit]
Description=Xitoy WebApp Backend (FastAPI)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/xitoy_webapp/backend
EnvironmentFile=/opt/xitoy_webapp/backend/.env
ExecStart=/opt/xitoy_webapp/backend/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
